# build/Debian/postinst
#!/bin/bash
chmod -R +rwx /opt/mybatop
systemctl daemon-reload
systemctl enable --now mybatop-shutdown.service
systemctl enable --now mybatop-startup.service
systemctl enable --now mybatop-status.service
ln -s /opt/mybatop/scripts/runscript/mybatop /usr/bin/mybatop
