# src/mybatop/filesystemd/mybatop-startup.service
[Unit]

Description="Run my shell script"

[Service]

ExecStart=/opt/mybatop/insert.sh

[Install]

WantedBy = multi-user.target
